# docker/webserver/supervisord.conf
[supervisord]
nodaemon = true

[program:nginx]
command = /usr/sbin/nginx
startsecs = 5
stdout_events_enabled = true
stderr_events_enabled = true

[program:app-gunicorn]
command = gunicorn -c "/etc/gunicorn.cfg" scripts/run_server:app 
stdout_events_enabled = true
stderr_events_enabled = true
